# root/etc/crontabs/root

0 1 * * * (echo "" && /bin/date && packt-cli -gd -c /config/configFile.cfg) 2>&1 | tee -a /config/logFile.log
